# docker/amazonlinux/updateexecution.ini
[program:updateexecution]
command=php /var/www/html/oil refine updateexecution
process_name=%(program_name)s_%(process_num)02d ; プロセスの命名規則、プロセスが複数の場合は必須
numprocs=4 ; プロセス数
autostart=true  ; supervisorが起動したら自動的に起動する
autorestart=true ; 落ちても自動的に再起動する
user=root
stopwaitsecs=600 ; supervisorからstop命令がされてから停止するまでの間の猶予時間。時間内に処理が終了しないと強制終了されてしまうため長めに設定する